# conf/php7-fpm.sv.conf
[program:php7-fpm]
command=php-fpm7.0 -F
